# confs/off1-reverse-proxy/systemd/system/certbot.service.d/override.conf
[Unit]
OnFailure=email-failures@certbot.service
