# confs/off-redis/systemd/system/redis.service.d/override.conf
[Unit]
# email on failures
OnFailure=email-failures@redis.service
