# confs/matomo/systemd/mariadb.service.d/override.conf
[Unit]
OnFailure=email-failures@mariadb-%l.service

