# confs/off2-reverse-proxy/systemd/system/nginx.service.d/override.conf
[Unit]
OnFailure=email-failures@nginx.service
Restart=always