# confs/ovh1-reverse-proxy/systemd/system/stunnel.target
[Unit]
Description=TLS tunnels for network services - per-config-file target

[Install]
WantedBy=multi-user.target
