# confs/off2-reverse-proxy/systemd/system/prometheus-nginx-exporter.service.d/override.conf
[Service]
# ONLY listen on 127.0.0.1 for security reasons
Environment="LISTEN_ADDRESS=127.0.0.1:9113"