# confs/moji-docker-prod-2/systemd/system/robotoff_db_backup.timer
[Unit]
Description=Launch Robotoff PostgreSQL backup every day

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
