# confs/off2-reverse-proxy/systemd/system/logrotate.service.d/override.conf
[Unit]
OnFailure=email-failures@logrotate.service
