# confs/matomo/systemd/redis-server.service.d/override.conf
[Unit]
OnFailure=email-failures@%H__redis-server.service

