# confs/matomo/systemd/redis-server.service.d/override.conf
[Unit]
OnFailure=email-failures@redis.service

