# confs/matomo/systemd/php7.3-fpm.service.d/override.conf
[Unit]
OnFailure=email-failures@php7.3-fpm-%H.service

