# confs/ovh3/systemd/system/nfs-server.service.d/override.conf
[Unit]
OnFailure=email-failures@nfs-server.service

