# confs/ovh-stunnel-client/systemd/system/logrotate.service.d/override.conf
[Unit]
OnFailure=email-failures@logrotate-%l.service
