# confs/off-stunnel-client/systemd/system/stunnel@.service.d/override.conf
[Unit]
# email on failures
OnFailure=email-failures@nginx.service

[Service]
# we need to enable putting pid file in runtime directory, with right permissions
# while still starting as root (needed by stunnel)
Group=stunnel4
RuntimeDirectory=stunnel-%i
RuntimeDirectoryMode=0775
