# confs/common/systemd/system/sanoid_check.timer
[Unit]
Description=Run Sanoid check Every 2 hours

[Timer]
# every two hours
OnBootSec=2h
OnUnitInactiveSec=2h

[Install]
WantedBy=timers.target
