# confs/off2/systemd/sanoid.service.d/override.conf
[Unit]
OnFailure=email-failures@sanoid-%l.service
