# confs/ovh3/systemd/sanoid.service.d/override.conf
[Unit]
# %l specifier is not available in systemd version 241 used on ovh3, use %i instead
OnFailure=email-failures@sanoid-%i.service
